# root/etc/services.d/unicorn-loadbalancer/run
#!/usr/bin/with-contenv bash
npm start --prefix=/unicorn
